# examples/scheduled-tasks/systemd/claude-watch-index-refresh.timer
[Unit]
Description=claude-watch example: fire index refresh every 6 hours
Documentation=https://github.com/hndrewaall/claude-watch

[Timer]
# Fire at 00:00, 06:00, 12:00, 18:00 user-local time. systemd.time(7) covers
# the OnCalendar grammar; change this string to your own cadence.
OnCalendar=00,06,12,18:00:00

# Also fire 5 minutes after boot / user-manager start, so a freshly-rebooted
# host catches up on the schedule without waiting for the next calendar slot.
OnBootSec=5min
OnUnitActiveSec=6h

# Run a missed tick on the next opportunity (host was asleep / laptop closed).
Persistent=true

# Spread firings across a small window so multiple identical timers on the
# same host don't all stampede at once. Tweak or remove freely.
RandomizedDelaySec=30s

# The .service unit this timer drives must share the base name with this
# .timer (claude-watch-index-refresh.service). Override here only if you've
# deliberately diverged.
Unit=claude-watch-index-refresh.service

[Install]
WantedBy=timers.target
