# container/cron.d/cw-default
# claude-watch default cron entries.
#
# Baked into the image at /etc/cron.d/cw-default. These entries run
# inside every container deployment by default — they don't reference
# host-specific paths, host-specific tools, or operator-specific
# state. Anything operator-specific (PR-watch loops, custom event
# emitters, etc.) belongs in /etc/cron.d/private, which is reserved
# for operator bind-mounts (see container/compose.yml).
#
# Why bake these specifically: every claude-watch deployment wants
# active-agents.json (the queue-minisite and friends consume it) and
# the metrics emit (work-queue-exporter scrapes the Prom-style
# textfile). Without these baked, every host needed separate cron
# setup — painful on macOS Docker Desktop especially (LaunchAgents,
# Full Disk Access prompts).
#
# Format: m h dom mon dow user command
#   - user MUST be `hndrewaall` (uid 1000); the cron daemon runs as
#     root but setuid's into the user named in each crontab line.
#   - claude-watch is on PATH for uid 1000 (image-baked
#     /usr/local/bin/claude-watch).
#   - State file path /var/lib/claude-watch/active-agents.json
#     matches the in-container path that's backed by the
#     `claude-watch-state` named volume (or the CW_STATE_PATH host
#     bind-mount, when set) — same path either way.
#   - stderr from claude-watch goes to /tmp/claude-watch-cron-*.log
#     so operators can inspect failures via
#     `docker compose exec <c> cat /tmp/claude-watch-cron-*.log`.
#
# DISABLING AN ENTRY (no sudo): each command below runs through the
#   `cw-cron-run <job-name> ...` guard. To stop/resume a job without
#   root, flip its flag file via `cw-cron-toggle disable|enable <name>`
#   (job-names: active-agents, metrics, stale-ready-check, queue-check,
#   watcher-health-check). Flags live under the uid-1000-writable,
#   named-volume-persistent /var/lib/claude-watch/cron-disabled. Takes
#   effect next tick; persists across redeploys. See docs/cron-toggle.md.

# Publish active-agents JSON state file (queue_id + alive per subagent)
# every minute. Consumers: work-queue-exporter (Prom has_live_owner
# gauge), queue-minisite (owner-status column), cron-queue-check.
* * * * *   hndrewaall /usr/local/bin/cw-cron-run active-agents /usr/local/bin/claude-watch active-agents --json --write-state /var/lib/claude-watch/active-agents.json > /dev/null 2>>/tmp/claude-watch-cron-active-agents.log

# Emit Prometheus-style metrics every minute. logger writes to syslog
# (or stderr inside the container — cron's -L 15 flag routes stderr to
# the supervisor log).
* * * * *   hndrewaall /usr/local/bin/cw-cron-run metrics /usr/local/bin/claude-watch metrics > /tmp/claude-watch-cron-metrics.log 2>&1

# Stale-ready queue-item watchdog. Every 5 minutes, scan
# `session-task queue ready --json` for pending items that have been
# ready past the threshold (default 6 minutes) and emit ONE aggregated
# `queue-stale-ready` claude-event per qualifying tick. Single-emit per
# queue id (ledger at /var/lib/claude-watch/stale-ready-state.json), so
# a long-stuck item only triggers once until it leaves the queue.
# Silent exit when nothing qualifies — no console spam.
*/5 * * * *  hndrewaall /usr/local/bin/cw-cron-run stale-ready-check /usr/local/bin/claude-watch stale-ready-check > /tmp/claude-watch-cron-stale-ready.log 2>&1

# Stuck / orphaned queue-item watchdog. Every 5 minutes, scan
# `session-task queue list --all --json` for ORPHANED items (running with a
# dead owning PID) and STUCK items (wedged, or running with a stale
# heartbeat) and emit `queue-orphaned` / `queue-stuck` claude-events. This
# is the IN-TREE equivalent of the out-of-tree Prometheus WorkQueueOrphaned
# / WorkQueueStuckSoft alert rules. Single-emit per (qid, condition); ledger
# at /var/lib/claude-watch/queue-check-state.json. EMISSION IS GATED behind
# `[queue_check] emit_events` in config.toml (DEFAULT FALSE) — this entry
# runs the detector every tick but stays silent until the toggle is flipped
# on (detection-only runs just log "emit_events is OFF" and exit 0).
*/5 * * * *  hndrewaall /usr/local/bin/cw-cron-run queue-check /usr/local/bin/claude-watch queue-check > /tmp/claude-watch-cron-queue-check.log 2>&1

# Event-watcher health check. Every minute, look for .json event
# files in the queue dir that have been sitting unconsumed for >2 min.
# If any exist, the event watcher is dead (container recreate, session
# clear, etc). Injects a [CLAUDE-WATCH] alert directly into the tmux
# pane via `claude-watch inject` — bypasses the event system entirely
# (which is the whole point: if events are accumulating, the event path
# is broken). The script now defaults its pane to `claude-watch inject`'s
# own resolution (so a non-container host can share it); the container
# pins the pane + its slash-command restart hint explicitly here to keep
# behavior identical to the pre-tools/ move.
#
# LAST-RESORT AUTO-RESTART (opt-in, ENABLED in-container here):
# CW_WATCHER_HEALTH_AUTORESTART=1 turns on the bounded self-heal fallback.
# When the inject (main-loop) path has fired CW_WATCHER_HEALTH_AUTORESTART_AFTER
# times (default 3, gated by the script's 600s cooldown ≈ 30 min of a
# sustained-blind bus) without the queue draining, the cron relaunches
# claude-event-watch itself, DETACHED. The watcher's fd-based flock singleton
# guard makes that a no-op whenever a real watcher is alive, so it can only win
# when the bus is genuinely unconsumed — this is the cure for the inject-only
# WATCHER(S) DOWN churn (events stuck >1.5h when no main loop acted on the
# inject). See the script header for the full cardinal-rule safety argument.
# The inject path STILL fires first and remains primary; this only kicks in
# after that path provably failed. Disable per-deployment by dropping the
# CW_WATCHER_HEALTH_AUTORESTART var (or set it to anything but "1").
* * * * *  hndrewaall CW_WATCHER_HEALTH_PANE=claude-container:0.0 CW_WATCHER_HEALTH_RESTART_HINT='Restart the event watcher via /claude-container:start-watchers' CW_WATCHER_HEALTH_AUTORESTART=1 /usr/local/bin/cw-cron-run watcher-health-check /usr/local/bin/cw-watcher-health-check > /dev/null 2>>/tmp/claude-watch-cron-watcher-health.log
