# examples/scheduled-tasks/launchd/org.gbre.claude-watch.index-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  claude-watch example LaunchAgent: periodic index refresh on macOS.

  This is the launchd parallel of systemd/claude-watch-index-refresh.{service,timer}.
  It runs every 6 hours under the logged-in user; no root required.

  The Label MUST match the plist filename minus .plist - that's how launchctl
  identifies agents. If you rename the file, rename the Label too.

  Customise:
    - StartInterval        cadence in seconds (21600 = 6 hours). Or swap for
                           StartCalendarInterval for clock-aligned firings.
    - EnvironmentVariables WHATEVER_REPO_DIR + INDEX_CMD as for the systemd unit.
    - StandardOutPath /
      StandardErrorPath    where to capture stdout/stderr. /tmp/ is fine for
                           an example; redirect to ~/Library/Logs/<name>.log
                           if you want them to survive reboot.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.gbre.claude-watch.index-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-lc</string>
        <string>cd "$WHATEVER_REPO_DIR" &amp;&amp; eval "$INDEX_CMD"</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>WHATEVER_REPO_DIR</key>
        <string>${HOME}</string>
        <key>INDEX_CMD</key>
        <string>echo "[claude-watch-index-refresh] tick at $(date -Iseconds) (replace INDEX_CMD with your real command)"</string>
    </dict>

    <!-- Fire every 21600 seconds (6 hours). -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <!-- Run a missed firing on the next opportunity (laptop was asleep). -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Don't keep alive: this is a oneshot, exit-on-completion job. -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/org.gbre.claude-watch.index-refresh.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/org.gbre.claude-watch.index-refresh.err</string>

    <!-- Throttle: don't relaunch more often than every 60 seconds even if the
         job exits immediately. -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
